# dwd_global_rad_api_server/services.d/dwd_global_rad_api_server/finish
#!/usr/bin/execlineb -S0

# Print an info message
echo "Application stopped."